# forge/script/sotd-publish.path
[Unit]
Description=Publish SOTD when forge/out/PUBLISH appears

[Path]
PathExists=%h/dox/st33v.com/forge/out/PUBLISH

[Install]
WantedBy=default.target

